# systemd/sidekick-dashboard.service
[Unit]
Description=Sidekick Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/sidekick/sidekick
ExecStart=/usr/bin/node src/dashboard.js
Restart=always
RestartSec=10
User=sidekick
Group=sidekick
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
