# systemd/sidekick-postgres.service
[Unit]
Description=Sidekick PostgreSQL (Docker)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
WorkingDirectory=/home/sidekick/sidekick/docker
ExecStart=/usr/bin/docker compose --env-file /home/sidekick/sidekick/.env -f /home/sidekick/sidekick/docker/docker-compose.yml up -d postgres
ExecStop=/usr/bin/docker compose --env-file /home/sidekick/sidekick/.env -f /home/sidekick/sidekick/docker/docker-compose.yml stop postgres
ExecStopPost=/usr/bin/docker compose --env-file /home/sidekick/sidekick/.env -f /home/sidekick/sidekick/docker/docker-compose.yml rm -f postgres
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
