# systemd/sidekick-influxdb.service
[Unit]
Description=Sidekick InfluxDB (Docker)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
WorkingDirectory=/home/sidekick/sidekick/docker
ExecStart=/usr/bin/docker compose -f /home/sidekick/sidekick/docker/docker-compose.yml up -d influxdb
ExecStop=/usr/bin/docker compose -f /home/sidekick/sidekick/docker/docker-compose.yml stop influxdb
ExecStopPost=/usr/bin/docker compose -f /home/sidekick/sidekick/docker/docker-compose.yml rm -f influxdb
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
